# scripts/host-setup/move-tg-to-nas/move-tg-to-nas.service
[Unit]
Description=Mueve descargas Telegram (local) al NAS segun categoria
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/move-tg-to-nas.sh
